# .config/systemd/user/backup.timer
[Unit]
Description=Backup periodically

[Timer]
OnBootSec=12h
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
