# .config/systemd/user/backup.service
[Unit]
Description=Backup periodically

[Service]
Type=oneshot
ExecStart=/home/tlegutko/.config/scripts/battery-monitor.sh
