# .config/systemd/user/battery-monitor.service
[Unit]
Description=Alarm when battery is low

[Service]
Type=oneshot
ExecStart=/home/tlegutko/.config/scripts/battery-monitor.sh
