# .config/etc/systemd/system/hosts.service
[Unit]
Description=Reset /etc/hosts file every 5min, so it's always back quickly.

[Service]
Type=oneshot
ExecStart=/home/tlegutko/.config/hosts.sh
